# install/97-synchighlights.rules
KERNEL=="wlan*", ACTION=="add", RUN+="/mnt/onboard/.adds/go-kobo-sync/exec_on_wifi_add.sh"
